# bin/guarani
#!/bin/bash
#
### BEGIN INIT INFO
# Provides:          guarani
# Required-Start:    $remote_fs $syslog
# Required-Stop:     $remote_fs $syslog
# Should-Start:      $network $time
# Should-Stop:       $network $time
# Default-Start:     2 3 4 5
# Default-Stop:      0 1 6
# Short-Description: Start and stop the guarani server daemon
# Description:       Controls the main Guarani server daemon "guarani"
#                    and its wrapper script "guarani".
### END INIT INFO
#
#set -e
#set -u
#${DEBIAN_SCRIPT_DEBUG:+ set -v -x}


#900150983cd24fb0d6963f7d28e17f72

user=signey
dirLog=/var/log/guarani
! test -e $dirLog && mkdir $dirLog && chown $user $dirLog


if [ "$USER" != "$user" ]; then
	su - $user -c "umask u=rwx,g=rwx,o=;exec $0 $@" &
	exit
fi


raiz=$(dirname $(dirname $(realpath $0)))
raiz="$(dirname $(cd "$(dirname "$0")";pwd))"
EU="$(cd "$(dirname "$0")"; pwd)/$(basename "$0")"

# carrega cpj=classes jar
# carrega maquina versao
x=$raiz/java/compConfig.sh
test -e $x && . $x

conf=$raiz/bin/loader.xml

#maq
sun="/usr/local/javag/bin/java"
sun=$javaDev
export JAVA_HOME=$(dirname $(dirname $sun))
echo "JAVA_HOME=$JAVA_HOME"
#read

#lixo - sunOp="-Djavax.net.ssl.trustStore" 
sunOp="-Xmx1024m"
sunOp=
sunClass=$raiz/classes/sun/loader
seg=1
porta=8080


#export LD_LIBRARY_PATH=/usr/local/xapian/lib:$LD_LIBRARY_PATH

#classpath
if [ 1 -eq 1 ]; then
	export CLASSPATH=$sunClass:$cpj
	#for i in $raiz/bin/pacotes/*.jar; do
	#	export CLASSPATH=$CLASSPATH:$i
	#done
	echo "===> $CLASSPATH"
	#read
else
	export CLASSPATH=$sunClass:$raiz/bin/pacotes/mail.jar
	export CLASSPATH=$sunClass:$raiz/bin/pacotes/mysql.jar
	export CLASSPATH=$CLASSPATH:$raiz/bin/pacotes/xapian/xapian_jni.jar:$raiz/bin/pacotes/lucene-core-3.4.0.jar
	export CLASSPATH=$CLASSPATH:$raiz/bin/pacotes/apache-math.jar
	export CLASSPATH=$sunClass:$raiz/bin/pacotes/sqlite-jdbc-3.19.3.jar
fi


#LANG=pt_BR.UTF-8
#export LANG=pt_BR.ISO-8859-1

. $raiz/bin/guarani.sh 

